# deploy/polymarket-bot.service.example
# Copy to /etc/systemd/system/polymarket-bot.service (adjust User + paths).
# sudo systemctl daemon-reload && sudo systemctl enable --now polymarket-bot
#
# Example below matches the simplest root-based setup used in the README walkthrough.
# If you create a non-root deploy user instead, change User/Group/WorkingDirectory/ExecStart.

[Unit]
Description=Polymarket copy-trading bot
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=root
Group=root
WorkingDirectory=/root/polymarket-copy-trading-bot
ExecStart=/root/polymarket-copy-trading-bot/.venv/bin/python main.py
Restart=always
RestartSec=15

# Hardening (optional; may need loosening if something breaks)
NoNewPrivileges=true

[Install]
WantedBy=multi-user.target
